# _release/99-modem-udev.rules
# Match SIMCom 7600 USB devices (AT / MBIM / ECM)
ATTRS{idVendor}=="1e0e", ATTRS{idProduct}=="9001", ACTION=="add", SUBSYSTEM=="tty", ENV{SIMCOM_MODEM}="1"
ATTRS{idVendor}=="1e0e", ATTRS{idProduct}=="9011", ACTION=="add", SUBSYSTEM=="tty", ENV{SIMCOM_MODEM}="1"
ATTRS{idVendor}=="1e0e", ATTRS{idProduct}=="9018", ACTION=="add", SUBSYSTEM=="tty", ENV{SIMCOM_MODEM}="1"

# Extra steps for ECM mode (PID 9018) — only if not auto-handled
ATTRS{idVendor}=="1e0e", ATTRS{idProduct}=="9018", \
  ACTION=="add", SUBSYSTEM=="tty", \
  RUN+="/sbin/modprobe option", \
  RUN+="/bin/sh -c 'echo 1e0e 9018 > /sys/bus/usb-serial/drivers/option1/new_id'"

# Create convenience symlinks for known interface numbers
ACTION=="add", SUBSYSTEM=="tty", ENV{SIMCOM_MODEM}=="1", ATTRS{bInterfaceNumber}=="02", SYMLINK+="UsbModemAT"
ACTION=="add", SUBSYSTEM=="tty", ENV{SIMCOM_MODEM}=="1", ATTRS{bInterfaceNumber}=="03", SYMLINK+="UsbModemAT2"
